# packaging/debian/postinst
#!/bin/bash
set -e

case "$1" in
  configure)
    # Create system user for the service
    if ! id coastal &>/dev/null; then
      useradd -r -m -s /bin/bash -d /opt/coastalclaw coastal
    fi

    # Data and workspace directories
    mkdir -p /var/lib/coastalclaw/data /var/lib/coastalclaw/workspace
    chown -R coastal:coastal /var/lib/coastalclaw /opt/coastalclaw

    # Enable and start services
    systemctl daemon-reload
    systemctl enable coastalclaw-server.service coastalclaw-daemon.service coastalclaw-web.service
    systemctl start coastalclaw-server.service coastalclaw-daemon.service coastalclaw-web.service || true

    echo ""
    echo "  CoastalClaw installed successfully."
    echo ""
    echo "  Web UI:  http://localhost:5173"
    echo "  API:     http://localhost:4747"
    echo ""
    echo "  Manage:  sudo systemctl status coastalclaw-server"
    echo "  Logs:    sudo journalctl -u coastalclaw-server -f"
    echo ""
    ;;
esac
