# packaging/debian/prerm
#!/bin/bash
set -e

case "$1" in
  remove|purge)
    systemctl stop coastalclaw-server.service coastalclaw-daemon.service 2>/dev/null || true
    systemctl disable coastalclaw-server.service coastalclaw-daemon.service 2>/dev/null || true
    ;;
esac
