# coastalos/systemd/coastal-architect.service
[Unit]
Description=CoastalClaw Self-Build Architect
After=coastal-server.service

[Service]
Type=oneshot
User=coastal
WorkingDirectory=/opt/coastalclaw
ExecStart=/usr/bin/node packages/architect/dist/index.js
Environment=CC_DATA_DIR=/var/lib/coastalclaw/data
Environment=CC_REPO_ROOT=/opt/coastalclaw
Environment=CC_OLLAMA_URL=http://127.0.0.1:11434
