# coastalos/systemd/coastal-vibevoice.service
[Unit]
Description=Coastal.AI VibeVoice ASR + TTS Service
After=network.target
ConditionPathExists=/opt/coastal-ai/coastalos/vibevoice/server.py

[Service]
Type=simple
User=coastal
WorkingDirectory=/opt/coastal-ai
Environment=VIBEVOICE_PORT=8001
Environment=VIBEVOICE_ASR_MODEL=microsoft/VibeVoice-ASR
Environment=VIBEVOICE_TTS_MODEL=microsoft/VibeVoice-Realtime-0.5B
ExecStart=/usr/bin/python3 coastalos/vibevoice/server.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
