# coastalos/systemd/coastal-daemon.service
[Unit]
Description=CoastalClaw Daemon (Scheduler + Voice)
After=coastal-server.service
Requires=coastal-server.service

[Service]
Type=simple
User=coastal
WorkingDirectory=/opt/coastalclaw
ExecStart=/usr/bin/node packages/daemon/dist/index.js
Restart=always
RestartSec=10
Environment=CC_SERVER_URL=http://localhost:4747
Environment=CC_DATA_DIR=/var/lib/coastalclaw/data

[Install]
WantedBy=multi-user.target
