# coastalos/systemd/coastal-architect.timer
[Unit]
Description=Run CoastalClaw Architect nightly at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=coastal-architect.service

[Install]
WantedBy=timers.target
