# coastalos/systemd/coastal-infinity.service
[Unit]
Description=CoastalClaw Infinity Vector Database
After=network.target

[Service]
Type=simple
User=coastal
WorkingDirectory=/var/lib/coastalclaw/data
ExecStart=/usr/local/bin/infinity --port 23817 --data-dir /var/lib/coastalclaw/infinity
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
