# coastalos/systemd/coastal-web.service
[Unit]
Description=Coastal.AI Web Portal (Vite static preview)
After=coastal-server.service
Wants=coastal-server.service

[Service]
Type=simple
User=coastal
WorkingDirectory=/opt/coastal-ai/packages/web
ExecStart=/usr/bin/node /opt/coastal-ai/node_modules/.bin/vite preview \
  --port 5173 \
  --host 127.0.0.1
Restart=always
RestartSec=5
Environment=VITE_CORE_PORT=4747

[Install]
WantedBy=multi-user.target
