# coastalos/systemd/coastal-vllm.service
[Unit]
Description=CoastalClaw vLLM Inference Server
After=network.target
ConditionPathExists=/usr/local/bin/vllm

[Service]
Type=simple
User=coastal
WorkingDirectory=/opt/coastalclaw
Environment=CC_VLLM_MODEL=llama3.2
ExecStart=/bin/sh -c '/usr/local/bin/vllm serve "$CC_VLLM_MODEL" --host 127.0.0.1 --port 8000'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
