# coastalos/systemd/coastal-airllm.service
[Unit]
Description=Coastal.AI AirLLM Inference Server (layer-streaming)
After=network.target
ConditionPathExists=/usr/local/bin/airllm-server

[Service]
Type=simple
User=coastal
WorkingDirectory=/opt/coastal-ai
Environment=CC_AIRLLM_MODEL=meta-llama/Meta-Llama-3-70B-Instruct
ExecStart=/usr/local/bin/airllm-server --model ${CC_AIRLLM_MODEL} --port 8002
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
